# etc/systemd/system/lidclose.service
[Unit]
Description=Sleep When Lid Closed

[Service]
Type=simple
ExecStart=/usr/share/lightdm/lidclose.sh

[Install]
