# etc/systemd/system/myunlock.service
[Unit]
Description=Sound When Unlocked

[Service]
Type=simple
ExecStart=/usr/share/lightdm/myunlock.sh

[Install]
